# crontab
#!/bin/sh

python -c 'from event.send_email import send_event_emails_to_subscribers; send_event_emails_to_subscribers()' >> /proc/$(cat /var/run/crond.pid)/fd/1 2>&1
